# actividad3/mi-servicio.service
[Unit]
Description=Servicio saludar

[Service]
ExecStart=/usr/local/bin/mi-servicio.sh

[Install]
WantedBy=multi-user.target
